# upstart/moloch-capture.conf
#/etc/init/moloch-viewer.conf
description "Moloch Capture"
author "mark"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on shutdown

respawn limit 10 5

script
    exec _TDIR_/bin/moloch-capture -c _TDIR_/etc/config.ini >> _TDIR_/logs/capture.log 2>> _TDIR_/logs/capture.err 
end script
